# addons/point_of_sale/tools/posbox/overwrite_after_init/etc/cron.daily/odoo
#!/bin/sh

wget 'http://localhost:8095/hw_drivers/check_certificate' >/dev/null